# etc/fmachine.service
[Unit]
Description=Fmachine

[Service]
ExecStart=/usr/local/bin/fmachine
User=fmachine
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
